# infra/production/notifier/systemd/frontman-notifier.service
[Unit]
Description=Frontman Production Notifier
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=deploy
Group=deploy
WorkingDirectory=/opt/frontman-notifier/current
EnvironmentFile=/opt/frontman-notifier/env
ExecStart=/opt/frontman-notifier/current/bin/frontman_notifier start
ExecStop=/opt/frontman-notifier/current/bin/frontman_notifier stop
Restart=on-failure
RestartSec=10

KillMode=process
TimeoutStopSec=30

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/frontman-notifier
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
